# config/offering.conf
[program:offering_operator]
directory=
command=
user=
autorestart=
redirect_stderr=
stdout_logfile=
